# startup_scripts/supervisord/looker.cfg

[program:looker]
command=/home/looker/looker/looker start
directory=/home/looker/looker
user=looker
numprocs=1
stdout_logfile=/var/log/looker/looker.log
stderr_logfile=/var/log/looker/looker.log
autostart=true
autorestart=true
startretries=100
startsecs=10
stopwaitsecs=10
stopsignal=TERM
priority=1
environment=HOME=/home/looker
